# projects/stargazer/inst/linux/etc/init.d/stargazer.gentoo.2007
#!/sbin/runscript

opts="reload"

DAEMON=/usr/sbin/stargazer
STARGAZER_OPTS=""
PIDFILE=/var/run/stargazer.pid

depend() {
	need net
}

start() {
	ebegin "Starting stargazer"
	start-stop-daemon --start --quiet --exec ${DAEMON} -- ${STARGAZER_OPTS}
	eend $?
}

stop() {
	ebegin "Stopping stargazer"
	start-stop-daemon --stop --quiet --pidfile ${PIDFILE} --retry=INT/60/KILL/5
	rm -f ${PIDFILE}
	eend $?
}

reload() {
	ebegin "Reloading stargazer rules"
	start-stop-daemon --stop --signal 1 --quiet --pidfile $PIDFILE
	return 0
	eend $?
}
